# scripts/openclaw.service
[Unit]
Description=OpenClaw AI Gateway
Documentation=https://docs.openclaw.ai
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/openclaw
ExecStart=/usr/bin/env node scripts/run-node.mjs gateway
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Pass environment
EnvironmentFile=-%h/.openclaw/.env
EnvironmentFile=-%h/openclaw/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw

# Hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
